# systemd/modinteractive.service
[Unit]
Description=modInteractive Camera/PIR Triggered HDMI Video Display
Documentation=https://github.com/WeAreTheArtMakers/modInteractive
After=graphical.target
Wants=graphical.target

[Service]
Type=simple
User=pi
Group=pi
SupplementaryGroups=audio video render input gpio
WorkingDirectory=/opt/modInteractive
ExecStart=/opt/modInteractive/venv/bin/python /opt/modInteractive/main.py
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=10
KillSignal=SIGINT
KillMode=control-group
Environment=PYTHONUNBUFFERED=1
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=PATH=/opt/modInteractive/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
