# systemd/modinteractive.service
[Unit]
Description=modInteractive Motion Triggered Video Display
Documentation=https://github.com/WeAreTheArtMakers/modInteractive
After=graphical.target
Wants=graphical.target

[Service]
Type=simple
WorkingDirectory=/opt/modInteractive
ExecStart=/opt/modInteractive/venv/bin/python /opt/modInteractive/main.py
Restart=always
RestartSec=5
User=pi
Group=pi
Environment=PYTHONUNBUFFERED=1
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target